# Freescale/sw/sdcard/template/etc/inittab
::sysinit:/etc/init.d/rcS

::respawn:/sbin/getty -L ttymxc1 115200 vt102

