# ARM/linux_qt_a9/sw/ramdisk/template/etc/inittab
::sysinit:/etc/init.d/rcS

::respawn:/sbin/getty -L 115200 ttyAMA0 vt102
::respawn:/sbin/getty -L 9600   tty1 vt102
::respawn:/sbin/getty -L 38400  tty2 vt102
::respawn:/sbin/getty -L 115200 tty4 vt102
::respawn:/sbin/getty -L 230400 tty5 vt102

