# ARM/linux_qt_a9/sw/ramdisk/template/etc/inittab
::sysinit:/etc/init.d/rcS

::respawn:/sbin/getty -L 9600 ttyAMA0 vt102

::respawn:/sbin/getty -L 9600    tty1 vt102 # LCD_CONSOLE



